# src/deploy/nginx_sup.conf
[program:nginx]
command = nginx -g "daemon off;"
stdout_logfile = /logs/nginx_daemon.log
redirect_stderr = true
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8